# conf/supervisor/apache2.conf
[program:apache2]
command=apache2-foreground
environment=USER="www-data"